# setup.timer
[Unit]
Description=timer for setup mcw

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
