# setup.service
[Unit]
Description=setup mcw
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup.sh

[Install]
WantedBy=multi-user.target
